# prerequisites/profile.d/01_xdg_init.sh
source "$(realpath $(dirname $0))/../xdg/home.sh"
source "$(realpath $(dirname $0))/../xdg/compliance.sh"
